# service/nfc_osc_client.service
[Unit]
Description=NfcOscClient

[Service]
ExecStart=/home/nfc/nfc-ocs-client/run_osc_nfc_client.sh

[Install]
WantedBy=multi-user.target
